# templates/setup_routes.service.j2
[Unit]
After=network.target

[Service]
ExecStart=/bin/sh /etc/ip_routes.sh

[Install]
WantedBy=default.target